# integration/snafu.service
[Unit]
Description=Snafu

[Service]
ExecStart=/bin/bash /opt/snafu-control.sh
ExecStop=/bin/sleep 2
